# spider/apps/spider-kpi/deploy/spider-kpi-refresh.service
[Unit]
Description=Spider KPI data refresh
After=network.target

[Service]
Type=oneshot
User=jpruit20
WorkingDirectory=/home/jpruit20/.openclaw/workspace/spider/apps/spider-kpi
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/jpruit20/.openclaw/workspace/spider/apps/spider-kpi/.env
ExecStart=/home/jpruit20/.openclaw/workspace/spider/apps/spider-kpi/.venv/bin/python /home/jpruit20/.openclaw/workspace/spider/apps/spider-kpi/scripts/refresh_all.py
