# spider/apps/spider-kpi/deploy/spider-kpi-daily-insights.service
[Unit]
Description=Spider KPI cross-source AI insight generation (Opus 4.7)
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/spiderclawkpi/spider/apps/spider-kpi
EnvironmentFile=/opt/spiderclawkpi/spider/apps/spider-kpi/.env
ExecStart=/opt/spiderclawkpi/spider/apps/spider-kpi/.venv/bin/python /opt/spiderclawkpi/spider/apps/spider-kpi/scripts/generate_daily_insights.py
StandardOutput=append:/var/log/spider-kpi-daily-insights.log
StandardError=append:/var/log/spider-kpi-daily-insights.log
