# spider/apps/spider-kpi/deploy/spider-kpi-refresh.timer
[Unit]
Description=Run Spider KPI refresh every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=spider-kpi-refresh.service
Persistent=true

[Install]
WantedBy=timers.target
